# systemd/knx.service
[Unit]
Description=EIB/KNX system
After=syslog.target

[Service]
Type=forking
ExecStart=/usr/bin/knxbackbone

[Install]
WantedBy=multi-user.target
